# dietpi/conf/dietpi-service
#!/bin/sh
### BEGIN INIT INFO
# Provides:          dietpi-service
# Required-Start:    $local_fs $remote_fs
# Should-Start:      $network ramlog
# Required-Stop:     $local_fs $network $remote_fs
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       DietPi boot shutdown reboot.
### END INIT INFO

#-- commands for update --
#update-rc.d -f dietpi-service remove
#rm /etc/init.d/dietpi-service
#cp /DietPi/dietpi/conf/dietpi-service /etc/init.d/dietpi-service
#chmod +x /etc/init.d/dietpi-service
#update-rc.d dietpi-service defaults 80 20

NAME="dietpi-service"

start(){
	#Copy DietPi to ramdisk
	if [ ! -f /DietPi/.ramdisk ]; then 
		echo -e "DietPi:"
		echo -e "Starting Ramdisk"
		#Copy DietPi globals/scripts to RAM
		/boot/dietpi/dietpi-ramdisk 0

		#Generate active flag
		echo -e "0" > /DietPi/.ramdisk
		
		#start bootup
		/DietPi/dietpi/boot
	fi
}

stop(){
	#Copy DietPi back to Disk 
	if [ -f /DietPi/.ramdisk ]; then
		#remove active flag.
		rm /DietPi/.ramdisk &> /dev/null
		
		#Copy DietPi globals/scripts back to DISK
		/DietPi/dietpi/dietpi-ramdisk 1
	fi
}

remove(){
    stop
    update-rc.d -f "$NAME" remove
    rm -fv "$0"
}

case "$1" in
	start)
		start
	;;
	stop)
		stop
	;;
	remove)
		remove
	;;
	restart)
		stop
		start
	;;
	*)
		echo "Usage: $0 {start|stop}"
	;;
esac