# rootfs/etc/udev/rules.d/10-wifi-disable-powermanagement.rules
KERNEL=="wlan*", ACTION=="add", RUN+="/sbin/iwconfig wlan0 power off"
