# rootfs/etc/udev/rules.d/10-wifi-disable-powermanagement.rules
SUBSYSTEM=="net" KERNEL=="wlan*", ACTION=="add", RUN+="/sbin/iwconfig %k power off"
