# rootfs/etc/systemd/system/cron.service.d/dietpi.conf
[Unit]
After=redis-server.service mariadb.service
