# rootfs/etc/cron.hourly/dietpi
#!/bin/bash
{
	#////////////////////////////////////
	# DietPi Cron.Hourly script
	#
	#////////////////////////////////////
	#
	# Info:
	# - Location /etc/cron.hourly/dietpi
	#
	#////////////////////////////////////

	#----------------------------------------------------------------
	# Main Loop
	#----------------------------------------------------------------
	# Time sync update, if mode 3 (hourly) is detected
	if grep -qi '^[[:blank:]]*CONFIG_NTP_MODE=3' /DietPi/dietpi.txt; then

		/DietPi/dietpi/func/run_ntpd 1 &> /dev/null &

	fi
	#----------------------------------------------------------------
	# Logging

	#----------------------------------------------------------------
	exit
	#----------------------------------------------------------------
}
