# rootfs/etc/systemd/system/dietpi-wifi-monitor.service
# Monitors loss of WiFi connection and automatically reconnects
[Unit]
Description=DietPi-WiFi-Monitor
After=network.target network-online.target

[Service]
Type=simple
ExecStart=/var/lib/dietpi/services/dietpi-wifi-monitor.sh

[Install]
WantedBy=multi-user.target
