# rootfs/etc/systemd/system/dietpi-ramlog.service
[Unit]
Description=DietPi-RAMlog
# Order 1
After=local-fs.target DietPi.mount var-log.mount -.mount
Before=syslog.service rsyslog.service

[Service]
Type=forking
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /var/tmp/dietpi/logs
ExecStart=/bin/dash -c '/boot/dietpi/func/dietpi-ramlog 0 2>&1 >> /var/tmp/dietpi/logs/dietpi-ramlog.log'
ExecStop=/bin/dash -c '/DietPi/dietpi/func/dietpi-ramlog 1 2>&1 >> /var/tmp/dietpi/logs/dietpi-ramlog.log'

[Install]
WantedBy=local-fs.target
