# rootfs/etc/systemd/system/dietpi-boot.service
[Unit]
Description=DietPi-Boot
# Order 3
Requires=dietpi-preboot.service
After=dietpi-preboot.service network.target networking.service
Before=getty-pre.target getty@tty1.service getty.target ssh.service dropbear.service

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=tty
ExecStart=/bin/bash -c '/DietPi/dietpi/boot 2>&1 | tee /tmp/dietpi-boot.log'

[Install]
WantedBy=multi-user.target
