# rootfs/etc/systemd/system/dietpi-fs_partition_resize.service
[Unit]
Description=DietPi-FS_partition_resize
#Order 0
DefaultDependencies=false
After=local-fs-pre.target
Before=local-fs.target swap.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/bash -c '/var/lib/dietpi/services/fs_partition_resize.sh | tee /var/tmp/dietpi/logs/fs_partition_resize.log'
StandardOutput=tty

[Install]
WantedBy=local-fs.target
Wants=local-fs-pre.target
