# rootfs/etc/systemd/system/dietpi-wifi-monitor.service
# Monitors loss of WiFi connection and automatically reconnects
[Unit]
Description=DietPi-WiFi-Monitor
Requires=dietpi-ramdisk.service
Wants=network-online.target
After=dietpi-boot.service network-online.target network.target

[Service]
ExecStartPre=/DietPi/dietpi/func/obtain_network_details
ExecStartPre=/bin/dash -c '[ -e /sys/class/net/wlan$(sed -n 2p /DietPi/dietpi/.network) ]'
ExecStart=/var/lib/dietpi/services/dietpi-wifi-monitor.sh

[Install]
WantedBy=multi-user.target
