# rootfs/etc/systemd/system/dietpi-postboot.service
[Unit]
Description=DietPi-PostBoot
# Order 3
Requisite=dietpi-preboot.service
After=dietpi-preboot.service

[Service]
Type=idle
RemainAfterExit=yes
StandardOutput=tty
ExecStart=/boot/dietpi/postboot

[Install]
WantedBy=multi-user.target
